# config/includes.chroot/etc/rc.local
#!/bin/bash
if [ -f /usr/bin/dconf ]; then
    echo "dconf database settings"
    sudo -u gdm dbus-launch dconf update
fi
exit 0
